# main/crontab
*/6 * * * * cd /app && python3 /app/src/hometamon.py
