# main/crontab
*/6 * * * * python3 /src/hometamon.py
# */6 * * * * python3 /src/hometask.py