# main/crontab
* * * * * python3 /metamon_code/main/src/hometamon.py 
