# main/crontab
*/7 * * * * python3 /src/hometamon.py
*/9 * * * * python3 /src/hometask.py