# main/crontab
* * * * * python3 -c "from /metamon_code/main/src import hometamon; hometamon.deploy_handler()"
