# main/crontab
*/10 * * * * python3 /metamon_code/main/src/hometamon.py --deploy 1 
