# main/crontab
* * * * * python /metamon_code/test/main.py
